# deploy/linux/lingbi-studio.service
[Unit]
Description=Lingbi Studio NotebookLM-like Workbench
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/lingbi-studio
EnvironmentFile=/opt/lingbi-studio/.env.production
Environment=BIND_HOST=127.0.0.1
ExecStart=/opt/lingbi-studio/start.sh
Restart=always
RestartSec=5
User=lingbi
Group=lingbi
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/lingbi-studio/.data /opt/lingbi-studio/logs

[Install]
WantedBy=multi-user.target
